# .crontab
10 * * * * offlineimap -u ttyui > /home/lu/.offlineimap.log
*  * * * * /home/lu/.bin/warn-if-low.sh
*  * * * * date --rfc-3339=seconds >> /home/lu/.pings
0  * * * * PYTHONPATH=/home/lu/k/news/libs/ python /home/lu/k/news/news.py 'http://rockpapershotgun.com' 'http://prog21.dadgum.com' 'http://codinghorror.com' > /home/lu/k/news/feed.html

* 0-7,19-23 * * * DISPLAY=:0 zenity --error --text "get off the computer!"
